# supervisord.conf
[supervisord]
loglevel=trace
user=root

[program:gateway]
command=./slice-gateway

[program:server]
command=./slice-n-dice-server
